# .docker/supervisord-mysqld.conf
[program:mysqld]
command=/app/.docker/start-mysqld.sh
numprocs=1
autostart=true
autorestart=true
