# .docker/supervisord-apache2.conf
[program:apache2]
command=/app/.docker/start-apache2.sh
numprocs=1
autostart=true
autorestart=true
